# scr/decompress_xrit.cron
#!/bin/sh
#
base=$(dirname $(realpath $0))

. /opt/etc/setup.sh
. $base/setup.sh
$base/decompress_xrit
